# tbot/config-samples/linux/etc/systemd/system/tbot-web.service
[Unit]
Description=T-Bot web

[Service]
Type=simple
ExecStart=/bin/bash /home/USER/tbot/script-tbot-web.sh

[Install]
WantedBy=multi-user.target
